# kathara_qemu/qemu-kathara-master/kvm/startUp.service
[Unit]
 Description=StartUp Katharà
 
 
[Service]
 Type=simple
 ExecStart=/usr/local/setup.sh

 
[Install]
 WantedBy=multi-user.target
